# test/cookbooks/test/templates/tool.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string><%= @tool_name %></string>
    <key>ProgramArguments</key>
    <array>
        <string><%= @binary_path %></string>
        <% @binary_arguments.each do |arg| %>
            <string><%= arg %></string>
        <% end %>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/<%= @tool_name %>.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/<%= @tool_name %>.err</string>
</dict>
</plist>
